# debian/evdaemon.upstart
description "Monitor Linux event devices and modify their behavior."
author "Tuomas Räsänen (tuos) <tuos@codegrove.org>"
version "0.1.1"

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon     
respawn

umask 0
kill timeout 5

exec evdaemon /dev/input/event3 /dev/input/event5
