# etc/init/evdaemon.conf.in
description "@PACKAGE_DESCRIPTION@"
author "@PACKAGE_AUTHOR@"
version "@VERSION@"

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon     
respawn

umask 0
kill timeout 5

exec evdaemon @EVDEV1@ @EVDEV2@
